# .cursor-server/data/User/History/-914e9b6/5GDf.service
[Unit]
Description=Christmas Claim Bot - Ultra Fast Version
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/itsbot
ExecStart=/usr/bin/python3 /root/itsbot/main_fast.py --offset 0
Restart=always
RestartSec=10
StandardOutput=append:/root/itsbot/bot_output.log
StandardError=append:/root/itsbot/bot_error.log

# Performance optimizations
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target

